# docker/allinone/supervisord.full.conf
[supervisord]
nodaemon=true
loglevel=info
pidfile=/tmp/supervisord.pid
logfile=/dev/stdout
logfile_maxbytes=0

[program:jaeger]
command=/usr/local/bin/jaeger --set receivers.otlp.protocols.http.endpoint=0.0.0.0:4318 --set receivers.otlp.protocols.grpc.endpoint=0.0.0.0:4317
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:serve]
command=/bin/sh -lc 'TRIAL="${DOJOZERO_TRIAL_SOURCE:-trial_sources/daily/nba.yaml trial_sources/daily/nfl.yaml}"; set -- dojo0 serve --host "${DOJOZERO_SERVE_HOST:-0.0.0.0}" --port "${DOJOZERO_SERVE_PORT:-8000}" --trace-backend jaeger --trace-ingest-endpoint "${DOJOZERO_TRACE_INGEST_ENDPOINT:-http://127.0.0.1:4318}"; for p in $TRIAL; do set -- "$@" --trial-source "$p"; done; exec "$@"'
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:arena]
command=/bin/sh -lc 'STATIC_ARGS=""; if [ -f /app/arena-static/index.html ]; then STATIC_ARGS="--static-dir /app/arena-static"; fi; exec dojo0 arena --host "${DOJOZERO_ARENA_HOST:-0.0.0.0}" --port "${DOJOZERO_ARENA_PORT:-3001}" --trace-backend "${DOJOZERO_ARENA_TRACE_BACKEND:-jaeger}" --trace-query-endpoint "${DOJOZERO_ARENA_TRACE_QUERY_ENDPOINT:-http://127.0.0.1:16686}" $STATIC_ARGS'
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
